# services/novaos_dashboard/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

; 1) Health server (must bind to $PORT for Render)
[program:health]
command=/bin/bash -lc 'uvicorn health:app --host 0.0.0.0 --port ${PORT:-8080}'
directory=.
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
priority=10

; 2) PerplexityFetcher agent (runs from REPO ROOT, not /app)
[program:perplexity_fetcher]
command=/usr/bin/env python3 -u agents/PerplexityFetcherAgent/entrypoint.py
directory=..
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
priority=20
