# services/novaos_dashboard/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
loglevel=info

[program:health]
command=/bin/bash -lc 'uvicorn health:app --host 0.0.0.0 --port ${PORT:-8080}'
directory=.
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=10
stopasgroup=true
killasgroup=true

[program:perplexityfetcheragent]
command=/usr/bin/env python3 -u agents/PerplexityFetcherAgent/entrypoint.py
directory=../..
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=20
stopasgroup=true
killasgroup=true

[program:novacoreagent]
command=/usr/bin/env python3 -u agents/NovaCoreAgent/entrypoint.py
directory=../..
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=25
stopasgroup=true
killasgroup=true

[program:actionexecutoragent]
command=/usr/bin/env python3 -u agents/ActionExecutorAgent/entrypoint.py
directory=../..
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=30
stopasgroup=true
killasgroup=true
