# services/novaos_dashboard/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
loglevel=info

; 1) Health server (must bind to $PORT for Render)
[program:health]
command=/bin/bash -lc 'uvicorn health:app --host 0.0.0.0 --port ${PORT:-8080}'
directory=.
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=10
stopasgroup=true
killasgroup=true

; 2) PerplexityFetcher agent
[program:perplexity_fetcher]
command=/usr/bin/env python3 -u agents/PerplexityFetcherAgent/entrypoint.py
directory=../..
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=20
stopasgroup=true
killasgroup=true

; 3) TrendAnalyzer agent (new)
[program:trend_analyzer]
command=/usr/bin/env python3 -u agents/TrendAnalyzerAgent/entrypoint.py
directory=../..
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=30
stopasgroup=true
killasgroup=true
